# luci-app-fan/root/etc/uci-defaults/luci-fan
#!/bin/sh

/usr/libexec/fan-control init

if ! uci -q get 'luci-fan.@luci-fan[0]' >/dev/null 2>&1; then
	uci -q add luci-fan luci-fan >/dev/null 2>&1
fi

uci -q get 'luci-fan.@luci-fan[0].mode' >/dev/null 2>&1 || uci -q set 'luci-fan.@luci-fan[0].mode=smart'
uci -q get 'luci-fan.@luci-fan[0].manual_pwm' >/dev/null 2>&1 || uci -q set 'luci-fan.@luci-fan[0].manual_pwm=70'
uci -q get 'luci-fan.@luci-fan[0].max_rpm' >/dev/null 2>&1 || uci -q set 'luci-fan.@luci-fan[0].max_rpm=3000'
uci -q get 'luci-fan.@luci-fan[0].poll_interval' >/dev/null 2>&1 || uci -q set 'luci-fan.@luci-fan[0].poll_interval=5'

# Only set on_temp if not already configured (avoid overwriting user settings)
current_on=$(uci -q get 'luci-fan.@luci-fan[0].on_temp' 2>/dev/null)
if [ -z "$current_on" ]; then
	uci -q set 'luci-fan.@luci-fan[0].on_temp=60'
fi

# Only set off_temp if not already configured
current_off=$(uci -q get 'luci-fan.@luci-fan[0].off_temp' 2>/dev/null)
if [ -z "$current_off" ]; then
	uci -q set 'luci-fan.@luci-fan[0].off_temp=30'
fi

uci -q commit luci-fan >/dev/null 2>&1

uci -q delete ucitrack.@luci-fan[-1]
uci -q commit ucitrack >/dev/null 2>&1

exit 0
